# supervisord.conf
[program:gecko_backend]
command= /server  # Replace with the actual path to your app binary 
autostart=true
autorestart=true
startsecs=1  # Wait 10 seconds before starting after a crash